# trunk/user/inadyn/inadyn-2.8.1/inadyn.service
[Unit]
Description=Internet Dynamic DNS Client
Documentation=man:inadyn
Documentation=man:inadyn.conf
Documentation=https://github.com/troglobit/inadyn
ConditionPathExists=/usr/etc/inadyn.conf
After=network-online.target
Requires=network-online.target

[Service]
Type=simple
ExecStart=/usr/sbin/inadyn -n -s

[Install]
WantedBy=multi-user.target
